# deploy/systemd/pgw-health.service
[Unit]
Description=PGW Health
After=network-online.target pgw-api.service
Wants=network-online.target

[Service]
User=pgw
Group=pgw
EnvironmentFile=/etc/pgw/pgw.env
ExecStart=/usr/local/bin/pgw-health
Restart=always
RestartSec=10s

# Resource Limits
LimitNOFILE=4096
MemoryMax=64M
TasksMax=512

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
