# deploy/systemd/pgw-agent.service
[Unit]
Description=PGW Agent
After=network-online.target pgw-api.service
Wants=network-online.target

[Service]
User=pgw
Group=pgw
EnvironmentFile=/etc/pgw/pgw.env
ExecStart=/usr/local/bin/pgw-agent
Restart=always
RestartSec=2s

# Resource Limits
LimitNOFILE=65536
MemoryMax=256M
TasksMax=2048

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
