# deploy/systemd/pgw-api.service
[Unit]
Description=PGW API
After=network-online.target
Wants=network-online.target

[Service]
User=pgw
Group=pgw
EnvironmentFile=/etc/pgw/pgw.env
ExecStart=/usr/local/bin/pgw-api
Restart=always
RestartSec=2s

# Resource Limits
LimitNOFILE=65536
MemoryMax=512M
TasksMax=4096

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
